# luicellas-stalker.timer
[Unit]
Description=Run Luicellas Stalker

[Timer]
OnBootSec=5min
OnCalendar=8..19:12/30:0
RandomizedDelaySec=1min

[Install]
WantedBy=timers.target
